# bin/ha_vig
#!/bin/bash
#
# haproxy
#
# chkconfig:   - 85 15
# description:  HAProxy is a free, very fast and reliable solution \
#               offering high availability, load balancing, and \
#               proxying for TCP and  HTTP-based applications
# processname: haproxy
# config:      /etc/haproxy/haproxy.cfg
# pidfile:     /var/run/haproxy.pid

# Source function library.
. /etc/rc.d/init.d/functions

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ "$NETWORKING" = "no" ] && exit 0

#determine DC  based on HALB configuration file
if [ ! -e /etc/HALB/halb.conf ]
then
 echo "/etc/HALB/halb.conf not found"
 exit 1
fi
      my_site=$(cat /etc/HALB/halb.conf  | grep DC | awk -F "=" '{print $2}')
      HAPLB_HOME=$(cat /etc/HALB/halb.conf  | grep HAPLB_HOME | awk -F "=" '{print $2}')

prog=$2
binary=$(echo $prog | cut -d. -f1)
exec="$HAPLB_HOME$my_site/bin/haproxy_$binary"
HA_CONF="$HAPLB_HOME$my_site/haproxy/$prog.cfg"

[ -e /etc/sysconfig/$prog ] && . /etc/sysconfig/$prog

lockfile=/var/lock/subsys/haproxy

check() {
    $exec -c -V -f $HA_CONF
}

start() {
    $exec -c -q -f $HA_CONF
    if [ $? -ne 0 ]; then
        echo "Errors in configuration file, check with $prog check."
        return 1
    fi

    echo -n $"Starting $prog: "
    # start it up here, usually something like "daemon $exec"
    daemon "$exec -D -f $HA_CONF -p /var/run/$prog.pid"
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
return $retval
}

stop() {
    echo -n $"Stopping $prog: "
    # stop it here, often "killproc $prog"
    killproc $prog
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    return $retval
}

restart() {
    $exec -c -q -f $HA_CONF
    if [ $? -ne 0 ]; then
        echo "Errors in configuration file, check with $prog check."
        return 1
    fi
    stop
    start
}

reload() {
    $exec -c -q -f $HA_CONF
    if [ $? -ne 0 ]; then
        echo "Errors in configuration file, check with $prog check."
        return 1
    fi
    echo -n $"Reloading $prog: "
    $exec -D -f $HA_CONF -p /var/run/$prog.pid -sf $(cat /var/run/$prog.pid)
    #$exec -p /var/run/$prog.pid -sf $(cat /var/run/$prog.pid)
    retval=$?
    echo
    return $retval
}

force_reload() {
    restart
}

fdr_status() {
    status $prog
}

case "$1" in
    start|stop|restart|reload)
        $1
        ;;
    force-reload)
        force_reload
        ;;
    check)
        check
        ;;
    status)
        fdr_status
        ;;
    condrestart|try-restart)
    [ ! -f $lockfile ] || restart
    ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|try-restart|reload|force-reload}"
        exit 2
esac
